# deployment/styleguide/webhook.template
#!/sbin/openrc-run

command=webhook
command_args="-hooks $PROJECT_ROOT/deployment/styleguide/hooks.json &"